# config/supervisord.conf
[program:webapp.com]
command = docker-compose -f /opt/project/docker-compose.yml up
stdout_logfile = /var/log/webapps/webapp.com.log
redirect_stderr = true
