# config/supervisord.conf
[program:webapp.com]
command = docker-compose -f /opt/project/docker-compose.yml up
